# supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=/path/to/your/api/command
autostart=true
autorestart=true

[program:task]
command=/path/to/your/task/command
autostart=true
autorestart=true